# rice/system-configs/fancontrol.service
[Unit]
Description=ThinkPad Fan Control Daemon
Documentation=https://github.com/YOUR_USERNAME/thinkpad-cyberpunk
After=multi-user.target

[Service]
Type=simple
# Update this path to match your installation directory
ExecStart=/home/sulaco/Projects/thinkpad-cyberpunk/utils/bin/fanctl daemon
Restart=on-failure
RestartSec=5s

# Security settings
# Run as root (required for fan control)
User=root
# Limit capabilities
CapabilityBoundingSet=
AmbientCapabilities=

# Process management
KillMode=process
KillSignal=SIGTERM

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fancontrol

[Install]
WantedBy=multi-user.target
