# rice/system-configs/fancontrol.service
[Unit]
Description=ThinkPad Fan Control Daemon
Documentation=https://github.com/adromero/cyber-thinkpad
After=multi-user.target

[Service]
Type=simple
# This path will be updated automatically during installation
ExecStart=__INSTALL_PATH__/utils/bin/fanctl daemon
Restart=on-failure
RestartSec=5s

# Security settings
# Run as root (required for fan control)
User=root
# Limit capabilities
CapabilityBoundingSet=
AmbientCapabilities=

# Process management
KillMode=process
KillSignal=SIGTERM

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fancontrol

[Install]
WantedBy=multi-user.target
